# cases/stability/reboot/reboot_test.service
[Unit]
Description=Reboot test

[Service]
ExecStart=/usr/bin/reboot_test.sh
KillMode=process
Type=idle

[Install]
WantedBy=multi-user.target
